# .codeploy/portable.fs/usr/lib/systemd/system/ops-es.service
[Unit]
Wants=es@ops\x2des.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c "echo provision ops-es Elasticsearch server"
